# install/puvefa_gunicorn.socket
[Unit]
Description=Gunicorn socket for puvefa

[Socket]
ListenStream=/run/gunicorn.sock

[Install]
WantedBy=sockets.target
